# deploy/arc.service
[Unit]
Description=Asterisk Response Classifier
After=network.target asterisk.service

[Service]
Type=simple
User=asterisk
Group=asterisk
WorkingDirectory=INSTALL_DIR
Environment=LD_LIBRARY_PATH=INSTALL_DIR/lib
ExecStart=INSTALL_DIR/arc -port 9099 -config INSTALL_DIR/config/phrases.yaml -model INSTALL_DIR/model
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
