# deploy/arc.service
[Unit]
Description=Asterisk Response Classifier (AEAP + Vosk)
After=network.target docker.service asterisk.service
Wants=asterisk.service

[Service]
Type=simple
User=asterisk
Group=asterisk
WorkingDirectory=/opt/asterisk-response-classifier
ExecStart=/opt/asterisk-response-classifier/bin/arc -port 9099 -config /opt/asterisk-response-classifier/config/phrases.yaml -vosk-url=ws://127.0.0.1:2700
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
