# deploy/arc.service
[Unit]
Description=arc AEAP bridge to speech-to-phrase
After=network.target speech-to-phrase.service
Requires=speech-to-phrase.service

[Service]
Type=simple
User=asterisk
Group=asterisk
WorkingDirectory=/opt/asterisk-response-classifier
ExecStart=/opt/asterisk-response-classifier/arc -port 9099 -config /opt/asterisk-response-classifier/config/ivr.yaml
Restart=on-failure

[Install]
WantedBy=multi-user.target
